# .config/systemd/user/borg.service
# [Unit]
Description=User backups

[Service]
Type=oneshot
ExecStart=~/dev/system/borg/backup.sh

