# .config/systemd/user/checkbattery.service
[Unit]
Description=Checking battery state

[Service]
Type=oneshot
ExecStart=/home/sheoak/bin/check-battery.sh
