# .config/systemd/user/checkbattery.service
[Unit]
Description=Checking battery state

[Service]
Type=oneshot
ExecStart=%h/bin/check-battery.sh
