# .config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service
#OnFailure=notify-failed@%i.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -aq -c ${DOTFILES_PRIVATE}/mbsyncrc
ExecStartPost=%h/.config/neomutt/scripts/post-sync
ExecStopPost=%h/.config/neomutt/scripts/post-sync
