# .config/systemd/user/mbsync.service.d/override.conf.sample
[Service]
Environment="HOME=/home/myusername/"
Environment="DOTFILES_PRIVATE=/home/myusername/.dotfiles-private/"
